# dnsalert.crontab
* 1 * * * /root/dnsalert/dnsalert-cron.bash
